# code/backup_daemon.service
[Unit]
Description=Simple Backup Daemon
After=network.target

[Service]
ExecStart=/usr/local/bin/backup_daemon.py
Restart=always
User=root
Group=root

[Install]
WantedBy=multi-user.target
